# sys/overlay/base/etc/cron.weekly/fstrim
#!/bin/sh

fstrim -Av --quiet-unsupported
